# crontab
# m h  dom mon dow   command
  0 1   *   *   *    /usr/local/rpi_telegraph/maintain.sh  > /usr/local/rpi_telegraph/update.log
